# docker/upstart/milon-actlg2.conf
description "milon www"
author "MP <mp@tcs.de>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Start the process
exec node /var/www/milonst/actionlogger2.js AMAZONTEST >> /mnt/milon-actionlogger2.log 2>&1