# docker/upstart/milon-st.conf
description "milon st"
author "MP <mp@tcs.de>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Start the process
exec node /var/www/milonst/server.js AMAZONTEST >> /mnt/milon-st.log 2>&1